# ansible/roles/common/files/etc/systemd/system/config-finished.service
[Unit]
Description=Config Finished Target
Requires=cloud-init.target
After=cloud-init.target

[Install]
WantedBy=multi-user.target
